# cmpt/opt/ceilometer/etc/systemd/system/ceilometer-agent-compute.service
[Unit]
Description=ceilometer-agent-compute

[Service]	
WorkingDirectory=/opt/ceilometer/bin
User=ceilometer
Group=ceilometer
ExecStart=/opt/ceilometer/bin/ceilometer-polling --polling-namespaces compute
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
